# docker/xray/templates/supervisord/daemon.ini
[program:x-ui]
command=x-ui
directory=/usr/local/bin
autorestart=true
priority=5
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=${LOGDIR}/supervisor/x-ui.out.log
stdout_logfile_maxbytes=0
stderr_logfile=${LOGDIR}/supervisor/x-ui.err.log
stderr_logfile_maxbytes=0
stopsignal=QUIT

[program:cron]
command=/usr/sbin/crond -f
autorestart=true
priority=10
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=${LOGDIR}/supervisor/cron.out.log
stdout_logfile_maxbytes=0
stderr_logfile=${LOGDIR}/supervisor/cron.err.log
stderr_logfile_maxbytes=0
stopsignal=QUIT

[program:dufs]
command=dufs -c ${WORKDIR}/dufs/conf.yml -a currycan:"${PASSWORD}"@/:rw
directory=/usr/local/bin
autorestart=true
priority=15
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=${LOGDIR}/supervisor/dufs.out.log
stdout_logfile_maxbytes=0
stderr_logfile=${LOGDIR}/supervisor/dufs.err.log
stderr_logfile_maxbytes=0
stopsignal=QUIT

[program:xray]
command=xray run -confdir ${WORKDIR}/xray/
directory=/usr/local/bin
autorestart=true
priority=20
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=${LOGDIR}/supervisor/xray.out.log
stdout_logfile_maxbytes=0
stderr_logfile=${LOGDIR}/supervisor/xray.err.log
stderr_logfile_maxbytes=0
stopsignal=QUIT

[program:v2ray]
command=v2ray run -c ${WORKDIR}/v2ray/config.json
directory=/usr/local/bin
autorestart=true
priority=20
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=${LOGDIR}/supervisor/v2ray.out.log
stdout_logfile_maxbytes=0
stderr_logfile=${LOGDIR}/supervisor/v2ray.err.log
stderr_logfile_maxbytes=0
stopsignal=QUIT

[program:nginx]
command=/usr/sbin/nginx -g "daemon off; error_log /dev/stderr info;"
autostart=true
autorestart=true
priority=25
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=${LOGDIR}/supervisor/nginx.out.log
stdout_logfile_maxbytes=0
stderr_logfile=${LOGDIR}/supervisor/nginx.err.log
stderr_logfile_maxbytes=0
stopsignal=QUIT
